# relay-package/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/dev.copper.relayserver.plist